# systemd/cognitord.service
[Unit]
Description=Cognitord - DSRs LLM Processing Daemon
Documentation=https://github.com/awdemos/cognitord
Requires=cognitord.socket

[Service]
Type=simple
User=cognitord
Group=cognitord
ExecStart=/usr/local/bin/cognitord --config /etc/cognitord/config.json
Restart=on-failure
RestartSec=5
StandardInput=socket
StandardOutput=socket
StandardError=journal
SyslogIdentifier=cognitord

# Environment variables
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/cognitord
ReadWritePaths=/var/log/cognitord
ReadWritePaths=/etc/cognitord
RemoveIPC=true

# Resource limits
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=80%

# Timeout settings
TimeoutStartSec=30
TimeoutStopSec=30

# Working directory
WorkingDirectory=/app

[Install]
WantedBy=multi-user.target